# etc/crontab
#!/bin/sh
* * * * * /usr/bin/php /var/www/html/vendor/bin/typo3 scheduler:run
