# etc/crontab
#!/bin/sh
* * * * * /usr/bin/php /var/www/html/public/typo3/sysext/core/bin/typo3 scheduler:run
